# whatsapp-backend.service
[Unit]
Description=WhatsApp Backend Server
After=network.target

[Service]
Type=simple
User=admin1
WorkingDirectory=/home/admin1/WhatsappAdvWebapp/whatsapp-advanced-webapp
ExecStart=/usr/bin/node server/whatsapp-server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3006
Environment=HOST=0.0.0.0
StandardOutput=append:/home/admin1/WhatsappAdvWebapp/whatsapp-advanced-webapp/logs/backend-service.log
StandardError=append:/home/admin1/WhatsappAdvWebapp/whatsapp-advanced-webapp/logs/backend-service-error.log

[Install]
WantedBy=multi-user.target
